# ubuntu_service/odroidu2-fan-service.conf
description "ODROID U2 FAN Service" 

start on runlevel [2345]

exec /bin/odroidu2-fan
